# docs/crontab
* * * * * /srv/jaguar/ENV/bin/python  /srv/jaguar/jaguarsite/log_loader.py
